# deploy/vps/vault/vault-unseal.service
# Auto-unseal Vault after Docker restarts
# Council adjustment #1: Vault must survive restarts without manual unseal
#
# Install:
#   cp vault-unseal.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable vault-unseal.service
#
# Security: init-keys.json must be chmod 600, owned by root
# This is a trade-off: convenience vs. security of storing unseal keys on disk
# For production HA, migrate to Vault auto-unseal (transit, AWS KMS, or GCP KMS)

[Unit]
Description=Vault Auto-Unseal
After=docker.service
Requires=docker.service
# Wait for the vault container to be running
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=/opt/vault/vault-unseal.sh
RemainAfterExit=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
