# ansible/roles/portal/templates/portal.service.j2
[Unit]
Description=webMethods Developer Portal
After=network.target

[Service]
Type=forking
User=root
Group=root
ExecStart=/opt/softwareag/portal/bin/startup.sh
ExecStop=/opt/softwareag/portal/bin/shutdown.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
