# deploy/vps/hegemon/hegemon-agent.service
[Unit]
Description=HEGEMON Agent — Claude Code tmux session
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hegemon
Environment=HOME=/home/hegemon
Environment=TERM=xterm-256color
Environment=LANG=en_US.UTF-8
ExecStart=/home/hegemon/.local/bin/hegemon-start.sh
ExecStop=/usr/bin/tmux -L hegemon kill-session -t hegemon
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
